# governors/ps5gov.service
[Unit]
Description=PS5 CPU/GPU load-adaptive governors
Documentation=file:/usr/local/lib/ps5-linux-cpuclock/governors/README.md
Conflicts=ps5boost.service ps5fan.service
After=multi-user.target dev-mp1.device
ConditionPathExists=/dev/mp1

[Service]
Type=simple
ExecStart=/usr/local/lib/ps5-linux-cpuclock/governors/run-governors.sh
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
